# middleware/spoolsense.service
[Unit]
Description=SpoolSense NFC Middleware
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/SpoolSense
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/SpoolSense/middleware/spoolsense.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
